# config/99-apple-touchbar.rules
# Prevent hid-sensor-hub from claiming Apple iBridge HID interfaces.
# The apple-ibridge driver needs ALL iBridge HID interfaces to function.
# Without this rule, hid-sensor-hub steals the display/brightness interface
# and the Touch Bar never activates.
ACTION=="add", SUBSYSTEM=="hid", ATTRS{idVendor}=="05ac", ATTRS{idProduct}=="8600", DRIVER=="hid-sensor-hub", RUN+="/bin/sh -c 'echo %k > /sys/bus/hid/drivers/hid-sensor-hub/unbind'"
